# Completed/Common_Core/Born2beroot/Changed_Files/monitoring.timer
[Unit]
Description=Run monitoring service every 10 minutes

[Timer]
OnBootSec=600s
OnUnitActiveSec=600s
Unit=monitoring.service

[Install]
WantedBy=timers.target
